# ops/webhook-deploy.service
[Unit]
Description=DPI-INTEROP Webhook Listener (GitHub → auto-deploy)
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/dpi-interop
EnvironmentFile=/etc/dpi-webhook.env
ExecStart=/usr/bin/python3 /opt/dpi-interop/ops/webhook-deploy.py
Restart=always
RestartSec=5

# Logs vers journald (en plus du LOG_FILE interne)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dpi-webhook

[Install]
WantedBy=multi-user.target
